# fitbit.crontab
0,15,30,45 * * * * /home/pace/fitbit/fitbit.cron
